# plgx-esp/default_data/default_queries/default_queries_windows_deep.conf
{
  "schedule": {
    "win_process_events": {
      "query": "select * from win_process_events ;",
      "interval": 30,
      "platform": "windows",
      "version": "2.9.0",
      "description": "Windows Process Events",
      "value": "Process Events",
      "removed": false
    },
    "win_file_events": {
      "query": "select * from win_file_events;",
      "interval": 180,
      "platform": "windows",
      "version": "2.9.0",
      "description": "File Integrity Monitoring",
      "value": "File Events",
      "removed": false
    },
    "win_removable_media_events": {
      "query": "select * from win_removable_media_events;",
      "interval": 240,
      "platform": "windows",
      "version": "2.9.0",
      "description": "Removable Media Events",
      "value": "Removable media Events",
      "removed": false
    },
    "win_epp_table": {
      "query": "select * from win_epp_table;",
      "interval": 360,
      "platform": "windows",
      "version": "2.9.0",
      "description": "Endpoint Products Status",
      "value": "Endpoint Products Status",
      "removed": false
    },
    "win_http_events": {
      "query": "select * from win_http_events;",
      "interval": 600,
      "platform": "windows",
      "version": "2.9.0",
      "description": "Windows HTTP Events",
      "value": "Windows HTTP Events",
      "removed": false
    },
    "win_socket_events": {
      "query": "select * from win_socket_events;",
      "interval": 240,
      "platform": "windows",
      "version": "2.9.0",
      "description": "Windows Socket Events",
      "value": "Windows Socket Events",
      "removed": false
    },
    "win_dns_events": {
      "query": "select * from win_dns_events;",
      "interval": 120,
      "platform": "windows",
      "version": "2.9.0",
      "description": "Windows DNS Events",
      "value": "DNS events",
      "removed": false
    },
    "win_registry_events": {
      "query": "select * from win_registry_events;",
      "interval": 300,
      "platform": "windows",
      "version": "2.9.0",
      "description": "Windows Registry Events",
      "value": "Registry events",
      "removed": false
    },
    "win_dns_response_events": {
      "query": "select * from win_dns_response_events;",
      "interval": 160,
      "platform": "windows",
      "version": "2.9.0",
      "description": "Windows DNS Response Events",
      "value": "DNS Response events",
      "removed": false
    },
    "win_ssl_events": {
      "query": "select * from win_ssl_events;",
      "interval": 160,
      "platform": "windows",
      "version": "2.9.0",
      "description": "Windows SSL Events",
      "value": "SSL connection events",
      "removed": false
    },
    "patches": {
      "query": "select hotfix_id, description, installed_on from patches;",
      "interval": 3600,
      "platform": "windows",
      "version": "2.9.0",
      "description": "Windows Patches",
      "value": "Installed Patches",
      "removed": false
    },
    "os_version": {
      "query": "select * from os_version;",
      "interval": 3600,
      "platform": "windows",
      "version": "2.9.0",
      "description": "os_version",
      "value": "OS version",
      "removed": false
    },
    "kernel_info": {
      "query": "select path, version, arguments from kernel_info;",
      "interval": 3600,
      "platform": "windows",
      "version": "2.9.0",
      "description": "Kernel Info",
      "value": "Kernel Info",
      "removed": false
    },
    "startup_items": {
      "query": "select name, path, status from win_startup_items;",
      "interval": 86400,
      "platform": "windows",
      "version": "2.9.0",
      "description": "Windows Startup Items",
      "value": "Startup Items",
      "removed": false
    },
    "drivers": {
      "query": "select * from drivers;",
      "interval": 86400,
      "platform": "windows",
      "version": "2.9.0",
      "description": "Windows Drivers",
      "value": "Installed drivers",
      "removed": false
    },
    "etc_hosts": {
      "query": "select * from etc_hosts;",
      "interval": 3600,
      "platform": "windows",
      "version": "2.9.0",
      "description": "Host file entries",
      "value": "Windows Host File Entries",
      "removed": false
    },
    "osquery_info": {
      "query": "SELECT * FROM osquery_info;",
      "interval": 28800,
      "description": "Information about the resident osquery process"
    },
    "wmi_cli_event_consumers": {
      "query": "SELECT * FROM wmi_cli_event_consumers;",
      "interval": 3600,
      "description": "WMI CommandLineEventConsumer, which can be used for persistence on Windows."
    },
    "wmi_script_event_consumers": {
      "query": "SELECT * FROM wmi_script_event_consumers;",
      "interval": 3600,
      "description": "WMI ActiveScriptEventConsumer, which can be used for persistence on Windows"
    },
    "users": {
      "query": "SELECT * FROM users;",
      "interval": 3600,
      "description": "Local system users."
    },
    "uptime": {
      "query": "SELECT * FROM uptime;",
      "interval": 3600,
      "description": "System uptime"
    },
    "certificates": {
      "query": "select common_name, issuer, self_signed, not_valid_after, path from certificates;",
      "interval": 3600,
      "description": "List all certificates in the trust store",
      "removed": false
    },
    "chrome_extensions": {
      "query": "select identifier, version, description, path from chrome_extensions;",
      "interval": 86400,
      "platform": "windows",
      "version": "2.9.0",
      "description": "Extensions in the Chrome browser",
      "value": "Chrome Extensions",
      "removed": false
    },
    "win_image_load_events": {
      "query": "select * from win_image_load_events;",
      "interval": 180,
      "platform": "windows",
      "version": "2.9.0",
      "description": "Extensions in the Chrome browser",
      "value": "Chrome Extensions",
      "removed": false
    },
    "ie_extensions": {
      "query": "select iee.path,iee.name,iee.version,(select sha1 from win_hash wh where wh.path=iee.path limit 1 ) as sha1  from ie_extensions iee;",
      "interval": 86400,
      "platform": "windows",
      "version": "2.9.0",
      "description": "Extensions in the Internet Explorer",
      "value": "IE Extensions",
      "removed": false
    },
    "scheduled_tasks": {
      "query": "select name,action,path,enabled,datetime(next_run_time,'unixepoch','UTC') from scheduled_tasks;",
      "interval": 3600,
      "platform": "windows",
      "version": "2.9.0",
      "description": "Windows scheduled_tasks",
      "value": "Scheduled tasks",
      "removed": false
    },
    "appcompat_shims": {
      "query": "select ach.*,(select sha1 from win_hash wh where wh.path=ach.path limit 1 ) as sha1 from appcompat_shims ach;",
      "interval": 3600,
      "platform": "windows",
      "version": "2.9.0",
      "description": "Windows scheduled_tasks",
      "value": "Scheduled tasks",
      "removed": false
    },
    "powershell_events_script_blocks": {
      "query": "select * from powershell_events;",
      "interval": 3600,
      "description": "Powershell script blocks reconstructed to their full script content, this table requires script block logging to be enabled. - ATT&CK T1086,T1064",
      "removed": false
    },
    "win_process_open_events": {
      "query": "select * from win_process_open_events;",
      "interval": 90,
      "platform": "windows",
      "version": "2.9.0",
      "description": "Windows Process Open Events",
      "value": "Windows Process Open Events",
      "removed": false
    },
    "win_remote_thread_events": {
      "query": "select * from win_remote_thread_events;",
      "interval": 90,
      "platform": "windows",
      "version": "2.9.0",
      "description": "Remote Thread Events",
      "value": "Remote Thread Events",
      "removed": false
    },
    "powershell_events": {
      "query": "select * from powershell_events;",
      "interval": 300,
      "platform": "windows",
      "version": "2.9.0",
      "description": "Power Shell Events",
      "value": "Power Shell Events",
      "removed": false
    }
  }
}